# configs/cron/crontab
20 4 * * * /scripts/update-osm.sh &>> /var/log/cron.log 2>&1
