# configs/cron/crontab
20 4 * * * /scripts/update-osm.sh &>> /data/update-osm_$(date +%s).log
